# alpine/package/apk-template/files/post-install.in
#!/bin/sh

# Optional OpenSSH setup: run only if sshd is available
if command -v sshd >/dev/null 2>&1 || [ -x /usr/sbin/sshd ]; then
  mkdir -p /root/.ssh /var/run/sshd
  # Drop-in configuration instead of editing main sshd_config
  mkdir -p /etc/ssh/sshd_config.d
  cat > "/etc/ssh/sshd_config.d/@PKGNAME@.conf" <<CONF
PermitRootLogin prohibit-password
PubkeyAuthentication yes
PasswordAuthentication no
ChallengeResponseAuthentication no
UsePAM no
AllowUsers root
CONF
  # Generate host keys if missing
  if [ ! -f /etc/ssh/ssh_host_rsa_key ]; then ssh-keygen -A; fi
fi

# App-specific directories: add here if needed (left empty in template)
@POST_INSTALL_EXTRA@

# Enable services if OpenRC is present
if command -v rc-update >/dev/null 2>&1; then
  if command -v sshd >/dev/null 2>&1 || [ -x /usr/sbin/sshd ]; then
    rc-update add sshd default 2>/dev/null || true
    rc-service sshd start 2>/dev/null || true
  fi
  rc-update add @PKGNAME@ default 2>/dev/null || true
  rc-service @PKGNAME@ start 2>/dev/null || true
fi

# Ensure runtime directory ownership
if id @PKGNAME@ >/dev/null 2>&1; then
  install -d -m 755 -o @PKGNAME@ -g dialout /var/lib/@PKGNAME@
fi

# Generate per-user SSH key for @PKGNAME@ if missing (requires ssh-keygen)
if id @PKGNAME@ >/dev/null 2>&1; then
  if command -v ssh-keygen >/dev/null 2>&1; then
    user="@PKGNAME@"
    home="$(getent passwd "$user" | cut -d: -f6 2>/dev/null)"
    [ -n "$home" ] || home="/var/lib/$user"
    install -d -m 700 -o "$user" -g "$user" "$home/.ssh"
    if [ ! -s "$home/.ssh/id_ed25519" ]; then
      # Create an ed25519 key without passphrase, comment indicates host
      su -s /bin/sh -c "ssh-keygen -q -t ed25519 -N '' -C '${user}@$(hostname)' -f '$home/.ssh/id_ed25519'" "$user" 2>/dev/null || true
    fi
  fi
fi

exit 0